# root/etc/services.d/opendkim/run
#!/command/with-contenv bash

echo "Starting opendkim"
/usr/sbin/opendkim -D -f -x /etc/opendkim/opendkim.conf
