# root/etc/services.d/postfix/run
#!/command/with-contenv bash

echo "Starting postfix"
/usr/sbin/postfix -c /etc/postfix start-fg


